# scripts/versifine-api.service
[Unit]
Description=Versifine API (Bun + Hono + Drizzle/Postgres)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=versifine
Group=versifine
WorkingDirectory=/opt/versifine/repo/apps/api
EnvironmentFile=/etc/versifine/api.env
# India-first: resolve "today"/"yesterday"/period queries on IST wall-clock so a
# late-night (00:00–05:30 IST) "chai today" doesn't land on the previous UTC day.
Environment=TZ=Asia/Kolkata
# Run the API straight from TypeScript source. Bun executes .ts natively, and
# running from source (not a bundle) avoids two classes of bundler bugs we hit:
#   - thrown AppErrors escaping to Bun's default text/plain 500, and
#   - string-concatenated dynamic imports (cycle-breakers) resolving relative
#     to a single bundled file and failing at runtime.
# The whole monorepo + node_modules are deployed, so source-run has everything.
ExecStart=/usr/local/bin/bun run src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
KillSignal=SIGINT

# Hardening
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/versifine

[Install]
WantedBy=multi-user.target
